# Example/Agent/agent.server
#!/bin/sh
#chkconfig: 2345 88 80
#description:agent
#processname:agent
########################################
DIR_SOFT=./
prog=agent
NAME_SCRIPT=run_daemon.py
########################################
. /etc/profile 
. /etc/rc.d/init.d/functions

function start(){
	cd  ${DIR_SOFT}
	if [ ! -x ${NAME_SCRIPT} ]
	then 
		chmod 777 ${NAME_SCRIPT}
	fi
	python ${NAME_SCRIPT} start 
	ret=$?
	if [ $ret -eq 0 ]; then
	    action $"Starting $prog: " /bin/true
	else
	    action $"Starting $prog: " /bin/false
	fi
	return 0
}


function stop(){
	cd  ${DIR_SOFT}
	if [ ! -x ${NAME_SCRIPT} ]
	then 
		chmod 777 ${NAME_SCRIPT}
	fi
	
	python ${NAME_SCRIPT} stop
	ret=$?
	if [ $ret -eq 0 ]; then
	    action $"Stopping $prog: " /bin/true
	else
	    action $"Stopping $prog: " /bin/false
	fi

}

restart(){
	cd  ${DIR_SOFT}
	python ${NAME_SCRIPT} restart
}

# See how we were called.
case "$1" in
	start)
		start
    	;;
  	stop)
    	stop
    	;;
  	status)
    	status
    	;;
  	restart)
  		restart
   	 	;;
  	*)
    echo $"Usage: $0 {start|stop|status}"
    exit 2
esac

exit $?
